# UniCron/MagicNet.cron
* 0 * * * /data/adb/modules/MagicNet/MagicNetD.sh